# conf/crontab
30 2 * * 1 root letsencrypt renew > /vpn/renew.log 2>&1 ; supervisorctl restart stunnel nghttpx 
# An empty line is required at the end of this file for a valid cron file.
